# fs/etc/profile.d/i3_terminal.sh
# needed for compositor compatibility
export TERMINAL=xfce4-terminal