# fs/etc/profile.d/polybar_vars.sh
set -a
source "${HOME}/.config/polybar/vars.txt"
set +a

