# fs/etc/profile.d/polybar_labels.sh
set -a
source "${HOME}/.config/polybar/labels.txt"
set +a

